# src/dist/rhel-linux-x86_64/systemd/groovy-webapp-library-search.service
[Unit]
Description=Groovy-based Web app Library Search
After=network.target

[Service]
Type=simple
User=webapp
EnvironmentFile=-/etc/sysconfig/groovy-webapp-library-search
ExecStart=/usr/bin/groovy-webapp-library-search
Restart=on-abort

[Install]
WantedBy=multi-user.target
