# src/dist/rhel-linux-x86_64/systemd/groovy-webapp-library-search-refresh.timer
[Unit]
Description=Restart Groovy web app Library Search service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
