# src/dist/rhel-linux-x86_64/systemd/malva-search-service-refresh.service
[Unit]
Description=Restart Malva Search Service

[Service]
Type=oneshot
ExecStartPre=/usr/bin/systemctl stop malva-search-service
ExecStart=/usr/bin/systemctl start malva-search-service
