# src/dist/rhel-linux-x86_64/systemd/malva-search-service@.service
[Unit]
Description=Malva Search Service %I
After=network.target

[Service]
Type=simple
User=malva
EnvironmentFile=-/etc/sysconfig/malva-search-service
ExecStart=/usr/bin/malva-search-service %i
Restart=on-abort

[Install]
WantedBy=multi-user.target
